# cleanup/snapshot_cleanup.service.j2
[Unit]
Description=Snapshot Cleanup Service

[Service]
Type=simple
User={{ user }}
Group={{ user }}
ExecStart={{ path_to_cleanup_script }}

[Install]
WantedBy=default.target