# state-sync/snapshot_state_sync.service.j2
[Unit]
Description=Snapshot State-Sync Service

[Service]
Type=simple
User={{ user }}
Group={{ user }}
ExecStart={{ path_to_state_sync_script }}

[Install]
WantedBy=default.target